# meta-lvgl-demo/recipes-graphics/lvgl/lvgl-demo-sdl2/lvgl-demo-sdl2.service
[Unit]
Description=Autostart LVGL demo

[Service]
#Type=oneshot
Type=simple
ExecStart=/usr/bin/lvgl-demo-sdl2

[Install]
WantedBy=multi-user.target
